# misc-config/setup-tablet/99-setup-tablet.rules
ACTION=="add", ATTRS{idVendor}=="056a", ATTRS{idProduct}=="0353", RUN+="/home/daisy/.setup-tablet/udev-script.sh"
